# .config/systemd/user/push_dotfiles.service
[Unit]
Description=Push Dotfiles to Github

[Service]
ExecStart=%h/.config/zsh/push_dotfiles.sh

[Install]
WantedBy=push_dotfiles.timer