# .config/systemd/user/push_dotfiles.timer
[Unit]
Description=Timer to push dotfiles to github

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
